# libmatic/scaffold/templates/scripts/launchd/com.libmatic.suggest-topics.weekly.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.libmatic.suggest-topics.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <!-- libmatic がインストールされた venv を呼ぶ。プロジェクト直下の uv 環境を想定。 -->
    <string>cd ~/projects/MY_PROJECT && uv run libmatic suggest-topics --config config/libmatic.yml</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <!-- 毎週月曜 09:00 JST -->
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/libmatic.suggest-topics.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/libmatic.suggest-topics.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
